# crontab
# HackerStack Comprehensive Crontab
# Posts to BOTH Twitter + Bluesky
# Engagement runs 20 minutes AFTER each post
# All times in UTC (EST=UTC-5)

# ========================================
# MONDAY
# 8AM(13:00), 10AM(15:00), 2PM(19:00), 3PM(20:00), 7PM(00:00)
# ========================================
0 13 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh post "â˜• Monday kickoff! What's your #1 goal this week? #IndieHacker #MondayMotivation"
20 13 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 15 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸš€ Monday momentum: One small win today = big progress. What's yours? #IndieHacker #BuildInPublic"
20 15 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 19 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ’¼ Monday afternoon power hour! What are you shipping? #IndieHacker"
20 19 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 20 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ… Monday wrap: What did you accomplish? Own it! #IndieHacker #MondayVibes"
20 20 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 0 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Monday night thoughts: Plan tomorrow. Sleep. Repeat. #IndieHacker #LateNight"
20 0 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# TUESDAY
# 9AM(14:00), 11AM(16:00), 5PM(22:00), 6PM(23:00), 8PM(01:00)
# ========================================
0 14 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ“ˆ Tuesday energy! What's the biggest challenge you're solving? #IndieHacker #TuesdayMotivation"
20 14 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 16 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ”§ Tuesday tools: Drop a tool that changed your workflow. I'll start: Cursor.sh #IndieHacker #Productivity"
20 16 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 22 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ† Tuesday evening: Reflect on wins. Small steps = big journeys. #IndieHacker #EveningVibes"
20 22 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 23 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ˜´ Tuesday night: Rest up. Tomorrow is another build day. #IndieHacker"
20 23 * * 2 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 1 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Late night building: Who's still up? Drop a ðŸ”¥ #IndieHacker #NightOwl"
20 1 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# WEDNESDAY
# 9AM(14:00), 12PM(17:00), 5PM(22:00), 7PM(00:00), 8:30PM(01:30)
# ========================================
0 14 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ’¡ Wednesday wisdom: The best time to start is now. What are you building? #IndieHacker #WednesdayMotivation"
20 14 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 17 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ„ Wednesday afternoon: Mid-week check-in. Stuck? Ask away! #IndieHacker #Help"
20 17 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 22 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒƒ Wednesday night: Ship something today. Any progress = good progress. #IndieHacker #ShipIt"
20 22 * * 3 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 0 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Hump day done! Celebrate the small wins. #IndieHacker #WednesdayNight"
20 0 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

30 1 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒŒ Late night dev session: Who's coding? Drop your stack! #IndieHacker #LateNightCoding"
50 1 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# THURSDAY
# 10AM(15:00), 12PM(17:00), 2PM(19:00), 6PM(23:00), 7PM(00:01), 9PM(02:00)
# ========================================
0 15 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ“š Thursday learning: What's one thing you learned this week? #IndieHacker #ThursdayMotivation"
20 15 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 17 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ› ï¸ Thursday tools deep-dive: Best AI tool you've discovered? #IndieHacker #Tools"
20 17 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 19 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh post "âš¡ Thursday afternoon power! One more push before Friday. #IndieHacker #Productivity"
20 19 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 23 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Thursday night: Reflect, plan, rest. What a week! #IndieHacker"
20 23 * * 4 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

1 0 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ›Œ Late Thursday: One more day until the weekend. Hang in there! #IndieHacker #AlmostFriday"
21 0 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 2 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ… 2AM built sessions: Who's still awake? Building in the quiet. #IndieHacker #NightOwl"
20 2 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# FRIDAY
# 8:30AM(13:30), 11AM(16:00), 1PM(18:00), 2PM(19:00), 7PM(00:00)
# ========================================
30 13 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh post "â˜• Friday morning: What's the plan? One week closer to your goals! #IndieHacker #FridayMotivation"
50 13 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 16 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸš€ Friday mid-day: Ship something before the weekend! What are you finishing? #IndieHacker #ShipItFriday"
20 16 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 18 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ• Friday afternoon: Almost there! One last push? #IndieHacker"
20 18 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 19 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ… Friday wrap-up: What shipped this week? Big or small, celebrate! #IndieHacker #FridayFeeling"
20 19 * * 5 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 0 * * 6 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Friday night: You earned it. Rest up for the weekend. #IndieHacker #FridayNight"
20 0 * * 6 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# SATURDAY
# 9AM(14:00), 11AM(16:00), 7PM(00:00), 9PM(02:00), 10PM(03:00)
# ========================================
0 14 * * 6 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ¿ Saturday vibes: Work smart, enjoy the weekend! What are you exploring? #IndieHacker #SaturdayVibes"
20 14 * * 6 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 16 * * 6 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ› ï¸ Saturday building: Who's working on a side project? Drop what you're building! #IndieHacker #SideProject"
20 16 * * 6 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 0 * * 7 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Saturday night: Reflect on the week. What went well? #IndieHacker #WeekendReflection"
20 0 * * 7 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 2 * * 7 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒŒ 2AM thoughts: Building something cool? Share! #IndieHacker #LateNight"
20 2 * * 7 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 3 * * 7 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ… 3AM builder check-in: The most productive hours. Who's with me? #IndieHacker #EarlyMorning"
20 3 * * 7 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# SUNDAY
# 8AM(13:00), 10AM(15:00), 6PM(23:00), 9PM(02:00)
# ========================================
0 13 * * 0 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ”‹ Sunday reset: One commitment for next week. What is it? #IndieHacker #SundayPrep"
20 13 * * 0 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 15 * * 0 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸ“ Sunday planning: Map out next week. Goals = achieved! #IndieHacker #SundayPlanning"
20 15 * * 0 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 23 * * 0 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ™ Sunday night: Prep your mind. Monday is coming. You've got this! #IndieHacker #SundayNight"
20 23 * * 0 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

0 2 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh post "ðŸŒ… 2AM Sunday/Monday boundary: Early birds, late night owls - same mission! #IndieHacker #Grind"
20 2 * * 1 /root/.openclaw/hackerstack/comprehensive-routine.sh engage

# ========================================
# ANALYTICS (twice daily)
# ========================================
0 8 * * * /root/.openclaw/hackerstack/gsc-auto.sh daily >> /root/.openclaw/hackerstack/cron/logs/gsc-daily.log 2>&1
0 20 * * * /root/.openclaw/hackerstack/gsc-auto.sh daily >> /root/.openclaw/hackerstack/cron/logs/gsc-daily.log 2>&1

# ========================================
# PRODUCT RESEARCH (every hour, 10 AM report)
# ========================================
0 * * * * /root/.openclaw/hackerstack/research-products.sh >> /root/.openclaw/hackerstack/cron/logs/research.log 2>&1

# Daily report at 10 AM EST (15:00 UTC)
0 15 * * * /root/.openclaw/hackerstack/daily-research-report.sh --send >> /root/.openclaw/hackerstack/cron/logs/daily-research.log 2>&1
